# ansible/roles/node-ollama-build/templates/node-ollama-build.timer.j2
[Unit]
Description=Weekly llama.cpp build timer

[Timer]
OnCalendar=weekly
RandomizedDelaySec=1800
Persistent=true

[Install]
WantedBy=timers.target
